# installer/launchd/dev.engram.brain.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Cold-brain LaunchAgent (optional).
  Install:
    cp installer/launchd/dev.engram.brain.plist.example \
       ~/Library/LaunchAgents/dev.engram.brain.plist
    # edit paths if needed, then:
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/dev.engram.brain.plist

  Runs every 2 hours (warm-ish). Uses pause-shell by default so Helix exclusive
  open is safe. Hot shell is dev.engram.local (always-on).
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.engram.brain</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-c</string>
    <string>set -a; source "$HOME/.engram/.env"; set +a; exec "$HOME/Engram/server/.venv/bin/engram" brain run --tier auto --profile quiet --pause-shell</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/REPLACE_ME/Engram/server</string>
  <key>StartInterval</key>
  <integer>7200</integer>
  <key>RunAtLoad</key>
  <false/>
  <key>Nice</key>
  <integer>10</integer>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>/Users/REPLACE_ME/.engram/logs/engram-brain.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/REPLACE_ME/.engram/logs/engram-brain.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/REPLACE_ME/Engram/server/.venv/bin:/Users/REPLACE_ME/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
